# embedded_boards/Raspberry_PI/Zero_2W/scripts_services/temp_logger_service/temp_logger.service
[Unit]
Description=CPU Temperature Logger

[Service]
ExecStart=/home/admin/temp_logger.sh
Restart=always
User=admin
WorkingDirectory=/home/admin

[Install]
WantedBy=multi-user.target